# supervisor/supervisord_sabnzbd.conf
[program:sabnzbd]
directory=/opt/sabnzbd
command=python SABnzbd.py
user = root
autostart=true
autorestart=true